# deploy/systemd/webex-headless-token-refresh.timer
[Unit]
Description=Refresh Webex Headless OAuth token cache periodically
[Timer]
OnBootSec=5m
OnUnitActiveSec=45m
RandomizedDelaySec=5m
Persistent=true
Unit=webex-headless-token-refresh.service

[Install]
WantedBy=timers.target
