# deploy/systemd/webex-headless-account-bot-sidecar-js.service
[Unit]
Description=Webex Headless JavaScript realtime sidecar for account bot
Documentation=https://github.com/Joey-Project/Webex-headless-messenger
After=network-online.target webex-headless-account-bot.service webex-headless-token-refresh.service
Wants=network-online.target webex-headless-token-refresh.service
Requires=webex-headless-account-bot.service
PartOf=webex-headless-account-bot.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=webex-headless-sidecar
Group=webex-headless-sidecar
WorkingDirectory=/opt/webex-headless-messenger/examples/sidecar-js
EnvironmentFile=/etc/webex-headless/webex-headless.env
ExecStartPre=/usr/bin/env WEBEX_SIDECAR_VALIDATE_CONFIG=1 WEBEX_SIDECAR_MESSAGE_EVENTS=created /usr/bin/node /opt/webex-headless-messenger/examples/sidecar-js/index.mjs
ExecStart=/usr/bin/env WEBEX_SIDECAR_MESSAGE_EVENTS=created /usr/bin/node /opt/webex-headless-messenger/examples/sidecar-js/index.mjs
Restart=on-failure
RestartSec=30s
KillSignal=SIGTERM
TimeoutStopSec=20s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
InaccessiblePaths=/etc/webex-headless/webex-headless-token.env
InaccessiblePaths=/var/lib/webex-headless-token
ReadOnlyPaths=/var/lib/webex-headless-access
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=webex-headless-account-bot.target
