# deploy/systemd/webex-headless-sidecar-receiver.service
[Unit]
Description=Webex Headless sidecar receiver
Documentation=https://github.com/Joey-Project/Webex-headless-messenger
After=network-online.target
Wants=network-online.target
PartOf=webex-headless-sidecar.target
StartLimitIntervalSec=300
StartLimitBurst=12

[Service]
Type=simple
User=webex-headless
Group=webex-headless
EnvironmentFile=/etc/webex-headless/webex-headless.env
ExecStart=/usr/local/bin/webex-headless sidecar receive --bind ${WEBEX_SIDECAR_RECEIVER_BIND} --path ${WEBEX_SIDECAR_RECEIVER_PATH} --health-path ${WEBEX_SIDECAR_RECEIVER_HEALTH_PATH}
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=20s
StateDirectory=webex-headless
ConfigurationDirectory=webex-headless
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/webex-headless
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=webex-headless-sidecar.target
