# deploy/systemd/webex-headless-sidecar-js.service
[Unit]
Description=Webex Headless JavaScript realtime sidecar
Documentation=https://github.com/Joey-Project/Webex-headless-messenger
After=network-online.target webex-headless-sidecar-receiver.service
Wants=network-online.target
Requires=webex-headless-sidecar-receiver.service
PartOf=webex-headless-sidecar.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=webex-headless
Group=webex-headless
WorkingDirectory=/opt/webex-headless-messenger/examples/sidecar-js
EnvironmentFile=/etc/webex-headless/webex-headless.env
ExecStartPre=-/usr/local/bin/webex-headless auth refresh --token-file ${WEBEX_TOKEN_FILE}
ExecStartPre=/usr/bin/env WEBEX_SIDECAR_VALIDATE_CONFIG=1 /usr/bin/node /opt/webex-headless-messenger/examples/sidecar-js/index.mjs
ExecStart=/usr/bin/node /opt/webex-headless-messenger/examples/sidecar-js/index.mjs
Restart=on-failure
RestartSec=30s
KillSignal=SIGTERM
TimeoutStopSec=20s
StateDirectory=webex-headless
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/webex-headless
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=webex-headless-sidecar.target
